# debian-package/etc/xdg/autostart/universal-battery-limiter.desktop
[Desktop Entry]
Name=Universal Battery Limiter
Comment=Battery charge limit management for all compatible laptops
Exec=/usr/local/bin/battery-indicator
Icon=battery
Terminal=false
Type=Application
Categories=System;Settings;
StartupNotify=true
X-GNOME-Autostart-enabled=true
X-GNOME-AutoRestart=true
X-Ubuntu-Gettext-Domain=universal-battery-limiter
Hidden=false
